# deploy/perr-bot.service
[Unit]
Description=PErr Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/PErr4.4.---
EnvironmentFile=-/root/PErr4.4.---/bot/.env
ExecStart=/root/PErr4.4.---/venv/bin/python /root/PErr4.4.---/bot/bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
