# deploy/perr-web.service
[Unit]
Description=PErr web chat (Flask + gunicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/PErr4.4.---/web
EnvironmentFile=-/root/PErr4.4.---/bot/.env
ExecStart=/root/PErr4.4.---/venv/bin/gunicorn -w 2 -b 127.0.0.1:5000 "app:app"
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
